# deploy/pre-earnings-agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.morConsulting.pre-earnings-agent</string>

    <key>ProgramArguments</key>
    <array>
        <!-- Apple Silicon: /opt/homebrew/bin/python3 -->
        <!-- Intel Mac:     /usr/local/bin/python3    -->
        <string>/opt/homebrew/bin/python3</string>
        <string>/Users/lukenmark/Documents/claude-os/tools/pre-earnings-agent/run.py</string>
        <string>start</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/lukenmark/Documents/claude-os/tools/pre-earnings-agent</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/lukenmark/Documents/claude-os/tools/pre-earnings-agent/logs/agent.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/lukenmark/Documents/claude-os/tools/pre-earnings-agent/logs/errors.log</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>30</integer>
</dict>
</plist>
